# root-image/etc/rc.local
LIVEUSER=bbqlinux

# Create the live user
useradd -m -k /etc/skel -g users -G audio,games,lp,nopasswdlogin,optical,power,scanner,storage,video,wheel -s /bin/bash -p $(mkpasswd --hash=md5 $LIVEUSER) $LIVEUSER
