# airootfs/etc/profile.d/editor.sh
# /etc/profile.d/editor.sh
export EDITOR=nano
