# root-image/etc/skel/.config/autostart/tint2.desktop
[Desktop Entry]
Type=Application
Exec=/usb/bin/tint2
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Tint 2
Comment=Taskbar
