# root-image/etc/systemd/system/create-liveuser.service
[Unit]
Description=Create the BBQLinux live user

[Service]
ExecStart=/etc/bbqlinux/create-liveuser

[Install]
WantedBy=multi-user.target 
